# etc/udev/rules.d/99-power-profile.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/bin/power-config.sh"
